# deploy/gitlab-runner/coder-templates-example/image/profile.d/00-coder-custom.sh
# Marker injected by the demo custom workspace image build (GitLab CI + Kaniko).
# Sourced by interactive login shells in any workspace based on this image.
echo "[coder-demo] custom workspace image (built by GitLab CI + Kaniko, air-gapped)"
